# infra/wiki-orchestrator.service
[Unit]
Description=Six Degrees of Wikipedia API Monthly Orchestrator Pipeline
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
# Assuming the repo is pulled somewhere like /opt/six-degrees-api
WorkingDirectory=/opt/six-degrees-api/infra

# Bun should be installed globally or referenced by absolute path
ExecStart=/usr/local/bin/bun run orchestrator.ts

# To ensure the process has permissions and doesn't get OOM-killed aggressively
OOMScoreAdjust=-500
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target